# supervisord.conf
[supervisord]
nodaemon=true

[program:sendmail]
command=service sendmail start

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"

[program:sethost]
command=/host_set.sh
